# buttplug/rootfs/etc/s6-overlay/s6-rc.d/buttplug/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: buttplug
# Runs the NodeJS daemon
# ==============================================================================

bashio::log.info "Starting Buttplug..."
exec npm start